# supervisord.conf
[supervisord]
nodaemon=true

[program:redis]
command=/usr/local/bin/redis-server

[program:encfs]
command=/usr/src/encKV/Proxy/TProxy 10
